# garage_script.conf
[program:garage_script]
command=python /home/pi/garage/server/server.py
autostart=true
autorestart=unexpected
stderr_logfile=/var/log/garage.err.log
stdout_logfile=/var/log/garage.out.log
